# record-and-playback/core/systemd/bbb-rap-events-worker.timer
[Unit]
Description=Run RaP events worker

[Timer]
OnBootSec=15min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
